# recipes-bsp/bootcount-reset/files/bootcount-reset.service
[Unit]
Description=Bootcount Reset timer

[Service]
Type=oneshot
ExecStart=/usr/bin/bootcount-reset

